# stubs/docker/app/supervisor.d/starter.ini
[program:starter]
process_name = %(program_name)s_%(process_num)02d
command = php -d variables_order=EGPCS %(ENV_APP_INSTALL_DIR)s/artisan toolbox:starter --mode=%(ENV_TOOLBOX_STARTER_MODE)s
autostart=true
autorestart=true
user = %(ENV_CURRENT_USER)s
loglevel = %(ENV_SUPERVISOR_LOG_LEVEL)s
redirect_stderr = true
stdout_logfile = %(ENV_APP_INSTALL_DIR)s/storage/logs/starter.log
stopwaitsecs=3600
numprocs=1
