# stubs/docker/app/supervisor.d/starter.ini
[program:starter]
directory = %(ENV_APP_INSTALL_DIR)s
command = php -d variables_order=EGPCS %(ENV_APP_INSTALL_DIR)s/artisan toolbox:starter --mode=%(ENV_TOOLBOX_STARTER_MODE)s

autostart = true
autorestart = true
startsecs = 1
startretries = 10
stopwaitsecs=60
stopsignal = TERM
stopasgroup = true
killasgroup = true

user = %(ENV_CURRENT_USER)s

loglevel = %(ENV_SUPERVISOR_LOG_LEVEL)s
redirect_stderr = true
stdout_logfile = %(ENV_APP_INSTALL_DIR)s/storage/logs/starter.log

priority = 1

numprocs=1
process_name = %(program_name)s
